# scripts/setup/touchscreen/systemd/zoe-kiosk-watchdog.service
[Unit]
Description=Zoe Kiosk Watchdog
After=zoe-kiosk.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zoe-kiosk-watchdog.sh
StandardOutput=journal
StandardError=journal
