# scripts/setup/touchscreen/systemd/zoe-kiosk.service
[Unit]
Description=Zoe Touch Panel Kiosk
After=graphical-session.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zoe
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/zoe/.Xauthority
ExecStart=/opt/TouchKio/start-kiosk.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
