# scripts/setup/jetson/zoe-graphify-refresh.service
[Unit]
Description=Refresh Zoe Graphify knowledge graph
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/zoe/assistant
ExecStart=/home/zoe/assistant/scripts/maintenance/refresh_graphify.sh --daily
StandardOutput=journal
StandardError=journal
TimeoutStartSec=7200
