# scripts/setup/jetson/zoe-graphify-refresh.timer
[Unit]
Description=Daily Zoe Graphify refresh

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
RandomizedDelaySec=30min
Unit=zoe-graphify-refresh.service

[Install]
WantedBy=timers.target
