# scripts/setup/jetson/zoe-multica-heartbeat.service
[Unit]
Description=Reflect Zoe runtime liveness into the Multica runtime registry
# Installed as a `systemctl --user` unit. System-scope ordering targets
# (docker.service, network-online.target) are not resolvable from the user
# manager, so we deliberately omit them: the script fails gracefully if the
# Multica DB container or Zoe API is briefly unavailable and the 60s timer
# simply retries on the next tick.

[Service]
Type=oneshot
WorkingDirectory=/home/zoe/assistant
ExecStart=/usr/bin/python3 /home/zoe/assistant/scripts/maintenance/multica_runtime_heartbeat.py
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60
