# scripts/setup/systemd/flue-zoe-telegram-watchdog.timer
# Runs the flue-zoe-telegram health watchdog once a minute (see the .service).
# OPTIONAL / operator opt-in. Enable with:
#   systemctl --user enable --now flue-zoe-telegram-watchdog.timer
#
# REQUIRES lingering (like the whole Zoe --user stack) so the timer keeps running
# when no session is logged in — otherwise the user manager (and the bot itself)
# stop at logout: `loginctl enable-linger $USER` (already enabled on Zoe hosts).

[Unit]
Description=Periodic health watchdog for flue-zoe-telegram

[Timer]
# First check shortly after boot (give the bot time to take the token over),
# then every minute. Persistent=false: a missed check while asleep is fine —
# the next tick catches an unhealthy bot within ~60s.
OnBootSec=2min
OnUnitActiveSec=60s
AccuracySec=15s

[Install]
WantedBy=timers.target
