# scripts/setup/systemd/llama-server.service
# llama.cpp server hosting Zoe's local voice brain — Gemma 4 E4B-QAT + MTP, port 11434.
#
# This mirrors the live Zoe config (the "brain" rock — see docs/CANONICAL.md). It is
# platform-specific; before installing, set:
#   - the llama-server binary path (built for your platform / GPU)
#   - the --model + --model-draft paths to your downloaded GGUFs
# See HARDWARE_COMPATIBILITY.md and scripts/setup/download_gguf_models.sh.
#
# HARD CONSTRAINTS (proven against the cached llama.cpp source — do not "fix" these):
#   - MTP speculative decoding (--spec-type draft-mtp) REQUIRES --flash-attn off on
#     this build (flash-attn on crashes in fattn.cu with the MTP drafter).
#   - With flash-attn OFF you may quantize the K cache (--cache-type-k q8_0) but the V
#     cache MUST stay f16 — a quantized V cache + FA off throws at startup
#     (llama-context.cpp). q8_0 on K only buys ~25-30% KV headroom so the ~3.8k-token
#     system prompt + tools fit in an 8192 context (a smaller ctx releases the slot on
#     overflow and the prompt KV is never cached → a cold prefill every turn).
[Unit]
Description=llama.cpp Server for Zoe (Gemma 4 E4B-QAT + MTP — live voice brain)
After=network.target

[Service]
Type=simple
LimitMEMLOCK=infinity

# ── Keep the brain rock resident ─────────────────────────────────────────────
# `--mlock` below asks llama.cpp to lock the weights; LimitMEMLOCK=infinity lets
# it. That is necessary but NOT sufficient — measured 2026-07-18 on the live
# Orin, VmLck held only 1.95 GB of a 5.6 GB RSS and 1,457 MB had still been
# swapped out (mlock covers the mapped model, not every CUDA/unified-memory
# allocation around it on Tegra).
#
# MemoryLow is soft (reclaim-resistant, not swap-proof), which is why it alone
# did not stop the eviction. MemorySwapMax=0 is the hard line for the pages
# mlock does not reach.
#
# MemoryLow was previously a LOCAL drop-in only (~/.config/systemd/user/
# llama-server.service.d/memory.conf) — real protection that existed on one box
# and in no template, so a rebuild or a second host silently lost it. Tracked
# here now. Same story for CPUWeight/IOWeight from priority.conf.
#
# NOTE: this is a ~5.6 GB RSS process. MemorySwapMax=0 removes its release
# valve, so pressure surfaces as an OOM kill instead of silent slow paging —
# which is the intended trade for the brain rock (a swapped brain is a broken
# brain), but it is why no MemoryMax is set here: a hard ceiling plus no swap
# would turn a transient spike into a kill. See docs/knowledge/voice-pipeline.md.
#
# No Nice=-N / OOMScoreAdjust=-N: a --user unit cannot raise priority on this
# box (`ulimit -e` is 0) and systemd SILENTLY DROPS the value while reporting
# success — verified 2026-07-18 (Nice=-5 applied as 0). Dev tooling
# de-prioritises itself instead (serena-mcp.service: Nice=10), which works
# because positive values need no privilege.
MemorySwapMax=0
MemoryLow=6G
CPUWeight=400
IOWeight=400
# The ExecStartPost health probe waits up to 120s for model+draft load; raise the
# start timeout above it so systemd doesn't kill the unit at the 90s default first.
TimeoutStartSec=180
# Point LD_LIBRARY_PATH at the directory of your llama.cpp build.
Environment=LD_LIBRARY_PATH=%h/llama.cpp/build-jetson-new/bin
# The endpoint is unauthenticated. Default to loopback; the live Zoe host binds
# 0.0.0.0 ONLY because it is fronted by Cloudflare Access + host firewalling. Do NOT
# expose 0.0.0.0 without equivalent auth/firewalling.
ExecStart=%h/llama.cpp/build-jetson-new/bin/llama-server \
  --model %h/models/gemma4-e4b-qat/gemma-4-E4B-it-qat-UD-Q4_K_XL.gguf \
  --model-draft %h/models/gemma4-e4b-qat/mtp-gemma-4-E4B-it.gguf \
  --spec-type draft-mtp \
  --spec-draft-n-max 4 \
  --spec-draft-p-min 0.6 \
  --spec-draft-ngl 99 \
  --host 127.0.0.1 \
  --port 11434 \
  --ctx-size 8192 \
  --cache-type-k q8_0 \
  --n-gpu-layers 99 \
  --parallel 1 \
  --cont-batching \
  --cache-reuse 256 \
  --flash-attn off \
  --temp 0.7 \
  --top-k 64 \
  --top-p 0.95 \
  --jinja \
  --chat-template-kwargs '{"enable_thinking":false}' \
  --mlock \
  --metrics
# Block until /health says ok (model + draft load), or fail the unit after 120s.
ExecStartPost=/bin/bash -c 'for i in $(seq 1 120); do curl -sf http://127.0.0.1:11434/health 2>/dev/null | grep -q ok && exit 0; sleep 1; done; exit 1'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
