# scripts/zoe-time-sync.service
[Unit]
Description=Zoe Time Synchronization Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/zoe
ExecStart=/usr/bin/python3 /home/pi/zoe/scripts/time_sync_service.py start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/pi/zoe/services/zoe-core

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/pi/zoe/data /home/pi/zoe/logs

[Install]
WantedBy=multi-user.target


