# services/zoe-data/systemd/bitnet-server.service
[Unit]
Description=BitNet b1.58 Inference Server (Pi ARM NEON)
Documentation=https://github.com/microsoft/BitNet
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/bitnet.cpp
# --model: path to downloaded BitNet 2B-4T gguf
# --host 127.0.0.1: localhost only, nginx proxies externally if needed
# --port 11434: OpenAI-compatible endpoint (matches BITNET_SERVER_URL in .env)
# --n-threads 4: Pi 5 has 4 performance cores
# -c 2048: context window (BitNet 2B is fast enough for this)
ExecStart=/home/pi/bitnet.cpp/build/bin/llama-server \
    --model /home/pi/models/bitnet/bitnet-b1.58-2B-4T-gguf/ggml-model-i2_s.gguf \
    --host 127.0.0.1 \
    --port 11434 \
    --n-threads 4 \
    --ctx-size 2048 \
    --log-disable
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
# Give the model 60s to load before health checks start
TimeoutStartSec=90

[Install]
WantedBy=default.target
