# services/zoe-data/systemd/gemma-server.service
[Unit]
Description=Gemma 4 E2B Inference Server (llama.cpp, Pi ARM NEON)
Documentation=https://github.com/ggml-org/llama.cpp
After=network.target bitnet-server.service
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/llama.cpp
# --model: Gemma 4 E2B Q4_K_M (~3.5 GB) from bartowski/google_gemma-4-E2B-it-GGUF
# Same model as Jetson — fine-tuned weights can be shared between devices
# --port 11435: separate from BitNet so both can run concurrently
# --threads 4: Pi 5 has 4 performance cores (use --threads not --n-threads for b8763+)
# --ctx-size 4096: Gemma 4 handles longer context for complex turns
ExecStart=/home/pi/llama.cpp/build/bin/llama-server \
    --model /home/pi/models/google_gemma-4-E2B-it-Q4_K_M.gguf \
    --host 127.0.0.1 \
    --port 11435 \
    --threads 4 \
    --ctx-size 4096 \
    --log-disable
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120

[Install]
WantedBy=default.target
