# services/zoe-data/systemd/mempalace.service
[Unit]
Description=MemPalace Hierarchical Memory Server (ChromaDB + SQLite)
Documentation=https://github.com/mem-palace/mempalace
After=network.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/pi
# MemPalace is a Python library (no server mode). This service just validates
# the installation and marks it as active for the zoe-data After= dependency.
# The Pi Agent uses MemPalace via direct Python API (from mempalace.searcher import search_memories).
ExecStart=/home/pi/.venv-zoe/bin/python -c "from mempalace.searcher import search_memories; print('MemPalace API ready')"
StandardOutput=journal
StandardError=journal
Environment=MEMPALACE_DATA_DIR=/home/pi/.mempalace

[Install]
WantedBy=default.target
