# config/archive/20251008_121422/zoe-dns-fix.service
[Unit]
Description=Zoe DNS Fix - Ensure zoe.local resolves correctly
After=network.target avahi-daemon.service
Wants=avahi-daemon.service

[Service]
Type=oneshot
ExecStart=/home/pi/fix-zoe-dns.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
